# etc/crontab
PATH=/srv/cdimage.ubuntu.com/bin:/usr/sbin:/usr/bin:/sbin:/bin

# On the Ubuntu production instance, this file lives in
# /srv/cdimage.ubuntu.com/etc/crontab; please update that file if you expect
# your changes to survive.

# CD image cron jobs, run out of cdimage's crontab.

# trusty daily builds
31 6 * * *	DIST=trusty for-project ubuntu cron.daily-live --live
14 4 * * *	DIST=trusty for-project kubuntu cron.daily-live --live
24 23 * * *	DIST=trusty for-project edubuntu cron.dvd --live
50 8 * * *	DIST=trusty for-project xubuntu cron.daily-live --live
29 5 * * *	DIST=trusty for-project ubuntu-server cron.daily --live
17 17 * * *	DIST=trusty for-project ubuntustudio cron.dvd --live
29 9 * * *	DIST=trusty for-project mythbuntu cron.daily-live --live
29 15 * * *	DIST=trusty for-project lubuntu cron.daily; DIST=trusty for-project lubuntu cron.daily-live --live
32 3 * * *	DIST=trusty for-project ubuntu-core cron.daily-live --live
14 21 * * *	DIST=trusty for-project ubuntukylin cron.daily-live --live
32 14 * * *	DIST=trusty for-project ubuntu-gnome cron.daily-live --live
# regular vivid daily builds
31 7 * * *	for-project ubuntu cron.daily-live --live
22 3 * * *	for-project ubuntu-desktop-next cron.daily-live --live
14 5 * * *	for-project kubuntu cron.daily-live --live
24 0 * * *	for-project edubuntu cron.dvd --live
50 9 * * *	for-project xubuntu cron.daily-live --live
29 6 * * *	for-project ubuntu-server cron.daily --live
17 18 * * *	for-project ubuntustudio cron.dvd --live
# mythbuntu will only release for LTS releases, but we keep daily builds around to ensure it remains working
29 10 * * *	for-project mythbuntu cron.daily-live --live
29 16 * * *	for-project lubuntu cron.daily; for-project lubuntu cron.daily-live --live
6 4 * * *	for-project livecd-base cron.daily-live --live
32 4 * * *	for-project ubuntu-core cron.daily-live --live
56 4 * * *	SUBPROJECT=system-image EXTRA_PPAS=snappy-dev/image for-project ubuntu-core cron.daily-preinstalled --live
# preinstalled armhf images
#35 1 * * *	for-project lubuntu cron.daily-preinstalled --live
14 22 * * *	for-project ubuntukylin cron.daily-live --live
32 15 * * *	for-project ubuntu-gnome cron.daily-live --live

# Also built on demand by touch release team.
02 2 * * *	for-project ubuntu-touch cron.daily-preinstalled --live
02 3 * * *	DIST=ubuntu-rtm/14.09 for-project ubuntu-touch cron.daily-preinstalled --live

*/15 * * * *	mirror-image-build-logs

# Self-rebuilds
*/5 * * * *	rebuild-requests -b -q trusty iso
*/5 * * * *	rebuild-requests -b -q vivid iso
