# etc/crontab
PATH=/srv/cdimage.ubuntu.com/bin:/usr/sbin:/usr/bin:/sbin:/bin

# On the Ubuntu production instance, this file lives in
# /srv/cdimage.ubuntu.com/etc/crontab; please update that file if you expect
# your changes to survive.

# CD image cron jobs, run out of cdimage's crontab.

# precise daily builds
31 8 * * *	DIST=precise for-project ubuntu cron.daily; DIST=precise for-project ubuntu cron.daily-live --live
29 9 * * *	DIST=precise for-project ubuntu-server cron.daily
45 9 * * *	DIST=precise for-project ubuntu-core cron.daily-live --live
56 21 * * *	DIST=precise SUBPROJECT=wubi for-project ubuntu cron.daily-live --live
34 14 * * 2,6	DIST=precise for-project ubuntu cron.dvd --live
19 19 * * *	DIST=precise for-project edubuntu cron.dvd --live
28 11 * * *	DIST=precise for-project kubuntu cron.daily; DIST=precise for-project kubuntu cron.daily-live --live
21 20 * * *	DIST=precise for-project mythbuntu cron.daily-live --live
34 20 * * *	DIST=precise for-project ubuntustudio cron.dvd --live
5 21 * * *	DIST=precise for-project xubuntu cron.daily; DIST=precise for-project xubuntu cron.daily-live --live
30 12 * * *	DIST=precise UBUNTU_DEFAULTS_LOCALE=zh_CN for-project ubuntu cron.daily-live --live
# regular raring daily builds
31 7 * * *	for-project ubuntu cron.daily-live --live
14 5 * * *	for-project kubuntu cron.daily-live --live
24 0 * * *	for-project edubuntu cron.dvd --live
50 9 * * *	for-project xubuntu cron.daily-live --live
29 6 * * *	for-project ubuntu-server cron.daily --live
17 18 * * *	for-project ubuntustudio cron.dvd --live
# mythbuntu will only release for LTS releases, but we keep daily builds around to ensure it remains working
29 10 * * *	for-project mythbuntu cron.daily-live --live
29 16 * * *	for-project lubuntu cron.daily; for-project lubuntu cron.daily-live --live
6 4 * * *	for-project livecd-base cron.daily-live --live
32 4 * * *	for-project ubuntu-core cron.daily-live --live
22 0 * * *	for-project kubuntu-active cron.daily-live --live
# preinstalled armhf images
35 1 * * *	for-project lubuntu cron.daily-preinstalled --live
32 13 * * *	for-project ubuntu cron.daily-preinstalled --live
# No Kubuntu Active on ARM until we know it works, jriddell 2012-02
#45 8 * * *	for-project kubuntu-active cron.daily-preinstalled --live
47 22 * * *	UBUNTU_DEFAULTS_LOCALE=zh_CN for-project ubuntu cron.daily-live --live
14 22 * * *	for-project ubuntukylin cron.daily-live --live
32 15 * * *	for-project ubuntu-gnome cron.daily-live --live

0 * * * *	mirror-image-build-logs
