# etc/crontab
PATH=/srv/cdimage.ubuntu.com/bin:/usr/sbin:/usr/bin:/sbin:/bin

# On the Ubuntu production instance, this file lives in
# /srv/cdimage.ubuntu.com/etc/crontab; please update that file if you expect
# your changes to survive.

# CD image cron jobs, run out of cdimage's crontab.

# xenial daily builds
41 6 * * *	DIST=xenial for-project ubuntu cron.daily-live --live
24 4 * * *	DIST=xenial for-project kubuntu cron.daily-live --live
00 1 * * *	DIST=xenial for-project xubuntu cron.daily-live --live
39 5 * * *	DIST=xenial for-project ubuntu-server cron.daily --live
27 17 * * *	DIST=xenial for-project ubuntustudio cron.dvd --live
39 9 * * *	DIST=xenial for-project mythbuntu cron.daily-live --live
39 15 * * *	DIST=xenial for-project lubuntu cron.daily; DIST=xenial for-project lubuntu cron.daily-live --live
42 3 * * *	DIST=xenial for-project ubuntu-base cron.daily-live --live
24 21 * * *	DIST=xenial for-project ubuntukylin cron.daily-live --live
42 14 * * *	DIST=xenial for-project ubuntu-gnome cron.daily-live --live
56 2 * * *	DIST=xenial for-project ubuntu-mate cron.daily-live --live
# regular yakkety daily builds
31 7 * * *	for-project ubuntu cron.daily-live --live
14 5 * * *	for-project kubuntu cron.daily-live --live
24 0 * * *	for-project edubuntu cron.dvd --live
50 1 * * *	for-project xubuntu cron.daily-live --live
29 6 * * *	for-project ubuntu-server cron.daily --live
17 18 * * *	for-project ubuntustudio cron.dvd --live
# mythbuntu will only release for LTS releases, but we keep daily builds around to ensure it remains working
29 10 * * *	for-project mythbuntu cron.daily-live --live
29 16 * * *	for-project lubuntu cron.daily; for-project lubuntu cron.daily-live --live
32 4 * * *	for-project ubuntu-base cron.daily-live --live
56 4 * * *	SUBPROJECT=system-image EXTRA_PPAS=snappy-dev/image for-project ubuntu-core cron.daily-preinstalled --live
#57 4 * * *	DIST=vivid SUBPROJECT=system-image EXTRA_PPAS=snappy-dev/image for-project ubuntu-core cron.daily-preinstalled --live
14 22 * * *	for-project ubuntukylin cron.daily-live --live
32 15 * * *	for-project ubuntu-gnome cron.daily-live --live
46 2 * * *	for-project ubuntu-mate cron.daily-live --live

# XXX: Note that the EXTRA_PPAS bit here also needs to be mirrored in rebuild-requests:
03 2 * * 2-6	DIST=vivid EXTRA_PPAS=ci-train-ppa-service/stable-phone-overlay:1001 for-project ubuntu-touch cron.daily-preinstalled --live
04 2 * * 2-6	DIST=xenial EXTRA_PPAS=ci-train-ppa-service/stable-phone-overlay:1001 for-project ubuntu-touch cron.daily-preinstalled --live
06 2 * * *	for-project ubuntu-touch cron.daily-preinstalled --live

*/15 * * * *	mirror-image-build-logs

# Self-rebuilds
*/5 * * * *	rebuild-requests -b -q vivid iso
*/5 * * * *	rebuild-requests -b -q xenial iso
*/5 * * * *	rebuild-requests -b -q yakkety iso
